# deploy/speckle/speckle-backup.timer
[Unit]
Description=Daily Speckle review server backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
